# supervisor/asyncfiles-worker.conf
[program:asyncfiles-worker]
; Команда запуска воркера
command=php %(ENV_DOCUMENT_ROOT)s/local/modules/adm.asyncfiles/bin/worker.php

; Количество воркеров (масштабирование)
numprocs=2
process_name=%(program_name)s_%(process_num)02d

; Автоматический перезапуск при падении
autostart=true
autorestart=true
startsecs=5
startretries=3

; Корректное завершение через SIGTERM (graceful shutdown)
stopsignal=SIGTERM
stopwaitsecs=30

; Логирование
stdout_logfile=/var/log/asyncfiles-worker-%(process_num)02d.log
stderr_logfile=/var/log/asyncfiles-worker-%(process_num)02d-error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

; Пользователь (bitrix-окружение)
user=www-data
